# archive-validator/systemd/synergy-archive-validator.service
[Unit]
Description=Synergy Archive Validator Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=synergy
Group=synergy
Environment=SYNERGY_ARCHIVE_CONFIG=/var/lib/synergy/archive-validator/config/archive-validator.toml
ExecStart=/usr/local/bin/synergy-archive serve
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
